# etc/init/xbmc-preload.conf
#!upstart
description "xbmc-preload"

start on filesystem

task

script
	start xbmc
end script